# infra/systemd/duckdns-refresh.service
[Unit]
Description=Refresh DuckDNS record for Brewgram
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/brewgram/duckdns.env
ExecStart=/usr/local/bin/update_duckdns.sh
